# rules.d/80-apa-test.rules
ACTION=="add", KERNEL=="sd?", SUBSYSTEM=="block", RUN+="apa-update-devices"
ACTION=="remove", KERNEL=="sd?", SUBSYSTEM=="block", RUN+=""
